# launchd/com.boyu.funding-agent.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.boyu.funding-agent.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/boyu/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>src.weekly_fetch</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/boyu/funding-agent</string>
    <key>StartCalendarInterval</key>
    <dict>
        <!-- Thursday = Weekday 4 -->
        <key>Weekday</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/boyu/funding-agent/outputs/logs/daily_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/boyu/funding-agent/outputs/logs/daily_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/boyu/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
