# launchd/com.boyu.funding-agent.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.boyu.funding-agent.backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>mkdir -p data/backups &amp;&amp; sqlite3 data/platform.db ".backup data/backups/platform-$(date +%Y%m%d).db"</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/boyu/funding-agent</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/boyu/funding-agent/outputs/logs/cli-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/boyu/funding-agent/outputs/logs/cli-backup.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/boyu/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
